# work/profile.d/email.sh
# Set EMAIL that git will use
export EMAIL="dylan.pinn@rea-group.com"
