# work/profile.d/rea-as.sh
# Use Google Authenticator for MFA method.
export REA_AS_MFA_METHOD=GOOGLE_AUTHENTICATOR
# Move config out of $HOME
export REA_AS_CACHE_PATH="${XDG_CACHE_HOME}/rea-as/"
