# work/profile.d/rea-as.sh
# Use Google Authenticator for MFA method.
export REA_AS_MFA_METHOD=GOOGLE_AUTHENTICATOR
