# work/profile.d/rea-as.sh
# Move rea-as cache out of $HOME.
export REA_AS_CACHE_PATH="$XDG_CACHE_HOME"/rea-as
# Use Google Authenticator for MFA method.
export REA_AS_MFA_METHOD=GOOGLE_AUTHENTICATOR
