# work/profile.d/update-machine.sh
update_machine () {
  brew update
  brew upgrade
  brew upgrade --cask
  rea --disable-auto-update cli update
}
