# mise/profile.d/default-npm-packages.sh
# shellcheck shell=sh

# Set a non-default location for default npm global packages.
export MISE_NODE_DEFAULT_PACKAGES_FILE="$XDG_CONFIG_HOME/mise/npm-default-packages"
