# homebrew/profile.d/homebrew.sh
# shellcheck shell=sh

# Setup Homebrew to work on any shell.
eval "$(/opt/homebrew/bin/brew shellenv)"
