# sh/profile.d/browser.sh
# Use lynx as text only TUI browser.
BROWSER=lynx
export BROWSER

