# sh/profile.d/pager.sh
# Use less as system pager.
PAGER=less
export PAGER
