# sh/profile.d/tz.sh
# Set timezone to be Melbourne.
TZ="Australia/Melbourne"
export TZ
