# sh/profile.d/visual.sh
# For visual editor use vim.
VISUAL='vim'
export VISUAL
