# sh/profile.d/visual.sh
# For visual editor use nvim (hopefully is available everywhere).
VISUAL='nvim'
export VISUAL
