# sh/profile.d/visual.sh
# shellcheck shell=sh

# Use vi as basic visual editor that should be installed most systems.
VISUAL="vi"
export VISUAL
