# sh/profile.d/visual.sh
# For visual editor use vim (hopefully is available everywhere).
VISUAL='vim'
export VISUAL
