# personal/profile.d/email.sh
# Set EMAIL that git will use
export EMAIL="me@dylanpinn.com"
