# personal/profile.d/homebrew-bundle.sh
# Use work Homebrew bundle file.
export HOMEBREW_BUNDLE_FILE=Brewfile-personal
